# first-cloud-boot/ubuntu-upstart.conf
start on stopped rc RUNLEVEL=[2345]

task

script
    if ! [ -e /root/xenserver.started ]; then
        touch /root/xenserver.started
        /bin/bash /root/boot-to-xenserver.sh
    fi
end script
